# viralpost.service
[Unit]
Description=ViralPost AI - Generador de Imágenes Virales
After=network.target

[Service]
Type=exec
WorkingDirectory=/home/user/AGT4
Environment="PATH=/home/user/AGT4/venv/bin"
EnvironmentFile=/home/user/AGT4/.env
ExecStart=/home/user/AGT4/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 5001 --workers 2
Restart=always
RestartSec=5

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=viralpost

[Install]
WantedBy=multi-user.target
